# crontab
0 0 * * *  %PROJECT_DIR%/console pullrequests:get >> %PROJECT_DIR%/logs/pullrequests_get.log